# mcbe_management/templete/mcbe.service
[Unit]
Description=Minecraft Bedrock Server Management Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/python/bin/python3 /var/games/mcbe/demon.py
ExecStop=/usr/local/python/binpython3 /var/games/mcbe/stop.py
Restart on-failure
User=root
Group=root
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target